# conf/ddupdate.cron
*/5 * * * * root /usr/local/bin/ddupdate
