# services/ddupdate.service
[Unit]
Description=ddupdate service

[Service]
Type=simple
ExecStart=/usr/local/bin/ddupdate

[Install]
WantedBy=default.target
